# jenkins/server-files/lib/systemd/system/jenkins.service.d/override.conf
[Service]
Environment="JENKINS_PORT=8071"